# dist/erlkoenig.service
[Unit]
Description=Erlkoenig Container Runtime
After=network.target

[Service]
Type=simple
User=root
Group=erlkoenig

Environment=ERL_EPMD_ADDRESS=127.0.0.1

# Wrapper loads cookie from /opt/erlkoenig/cookie, then delegates to relx.
# No RELX_COOKIE env var needed.
ExecStart=/opt/erlkoenig/bin/erlkoenig foreground

KillSignal=SIGTERM
TimeoutStopSec=30

Restart=on-failure
RestartSec=5

Delegate=yes
LimitNOFILE=65536
LimitMEMLOCK=infinity
NoNewPrivileges=no

[Install]
WantedBy=multi-user.target
